# deploy/fish-server.service
# systemd unit for the Fish game server.
#
# Install on the VM:
#   sudo cp deploy/fish-server.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now fish-server
#
# View logs:  journalctl -u fish-server -f

[Unit]
Description=Fish game server (Fastify + Socket.IO)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/fish
ExecStart=/usr/bin/node apps/server/dist/index.js
Restart=always
RestartSec=3
User=fish
Environment=NODE_ENV=production
EnvironmentFile=/opt/fish/apps/server/.env

[Install]
WantedBy=multi-user.target
